# systemd_worx_gps_rec.service
[Unit]
Description=Worx GPS Recorder Service
After=network-online.target
Wants=network-online.target
Documentation=file:///home/nilsgollub/Worx_GPS/RASPBERRY_PI_DEPLOYMENT.md

[Service]
Type=simple
User=nilsgollub
WorkingDirectory=/home/nilsgollub/Worx_GPS
ExecStart=/home/nilsgollub/Worx_GPS/.venv/bin/python3 /home/nilsgollub/Worx_GPS/Worx_GPS_Rec.py
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=worx_gps_rec

# Environment
Environment="PATH=/home/nilsgollub/Worx_GPS/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/nilsgollub/Worx_GPS/data /tmp

# Performance
CPUQuota=50%
MemoryMax=256M

[Install]
WantedBy=default.target
